# Solr 5.2/supervisord.conf
[supervisord]
nodaemon=true

[program:solrServer]
command=/solr/bin/solr "start" "-f"

