# Solr 5.2/supervisord.conf
[supervisord]
nodaemon=true

[program:solrServer]
command=/bin/bash -c "/solr/bin/solr start"

